# .config/bspwm/Desktops/bspwm/scripts/Reminders/reminders.timer
[Unit]
Description= run my reminders every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
